# .config/systemd/user/checkupdates.timer
[Unit]
Description=Run my update notification script every hour.

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target