# .config/systemd/user/checkupdates.timer
[Unit]
Description=Run my update notification script every hour.

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target